# supervisord.conf
[supervisord]
nodaemon=true

[program:yt]
command=python yt.py
autostart=true
autorestart=true
